# orchestrator/ai.myndaix.fix-sweep.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>ai.myndaix.fix-sweep</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/stevenfernandez/code/active/myndaix-runtime/orchestrator/fix-sweep.sh</string>
  </array>
  <key>StartInterval</key>
  <integer>3600</integer>
  <key>RunAtLoad</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/Users/stevenfernandez/.myndaix/orchestrator/fix-sweep.out</string>
  <key>StandardErrorPath</key>
  <string>/Users/stevenfernandez/.myndaix/orchestrator/fix-sweep.out</string>
</dict>
</plist>
